# _home_pi_.config_autostart/PhotoBooth.desktop
[Desktop Entry]
Name=PhotoBooth
Type=Application
Comment=Taking Pictures
Exec=/home/pi/42nd_PhotoBooth/loop.sh